# sanderson-rpg.service
[Unit]
Description=Sanderson RPG Character Management System
After=network.target

[Service]
Type=forking
User=gamebawx
WorkingDirectory=/home/gamebawx/gitrepos/sanderson_rpg/project_sanderson
ExecStart=/home/gamebawx/gitrepos/sanderson_rpg/project_sanderson/start.sh
ExecStop=/home/gamebawx/gitrepos/sanderson_rpg/project_sanderson/stop.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"
Environment="PORT=3000"

[Install]
WantedBy=multi-user.target
